# scripts/service/codero.service
[Unit]
Description=codero daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/codero daemon
Restart=on-failure
RestartSec=5
PIDFile=/var/run/codero/codero.pid

[Install]
WantedBy=multi-user.target
